# systemd/kdevops-slack-billing.service
[Unit]
Description=kdevops AWS Billing Slack Notification
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kdevops
Group=kdevops
WorkingDirectory=/data/cloud/gpt2/kdevops
Environment="KDEVOPS_DIR=/data/cloud/gpt2/kdevops"
ExecStart=/data/cloud/gpt2/kdevops/scripts/slack-billing-notify.sh

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/data/cloud/gpt2/kdevops

# Restart policy for failed notifications
Restart=on-failure
RestartSec=5min
StartLimitBurst=3
StartLimitIntervalSec=1h

[Install]
WantedBy=multi-user.target