# playbooks/roles/nix-cache-mirror/templates/nix-cache-sync.timer.j2
[Unit]
Description=Run Nix Binary Cache Sync Service hourly
Requires=nix-cache-sync.service

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
